# src/main/resources/org/jenkinsci/test/acceptance/docker/fixtures/JabberContainer/supervisord.conf
[supervisord]
nodaemon=true

[program:prosody]
user=prosody
command=/usr/bin/lua5.1 /usr/bin/prosody

[program:xmpp-client]
command=/usr/bin/python /usr/local/bin/logbot --host localhost --user bot@localhost --passwd bot-pw test@conference.localhost
stdout_logfile=/tmp/logbot.log