# src/main/resources/org/jenkinsci/test/acceptance/docker/fixtures/SvnContainer/config/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "exec /usr/sbin/apache2 -DFOREGROUND"

[program:svn]
command=svnserve -d -r /home/svn/

[program:ssh]
command=/usr/sbin/sshd -D -e