# src/main/resources/org/jenkinsci/test/acceptance/docker/fixtures/SvnContainer/config/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && /usr/sbin/apache2 -DFOREGROUND"

[program:svnserve]
command=/usr/bin/svnserve -d --foreground -r /var/lib/svn/myrepo/

[program:ssh]
command=/usr/sbin/sshd -D -e
