# src/main/resources/org/jenkinsci/test/acceptance/docker/fixtures/XMPPContainer/supervisord.conf
[supervisord]
nodaemon=true

[program:prosody]
#command=/usr/bin/prosodyctl start
command=/usr/bin/lua5.1 /usr/bin/prosody

[program:xmpp-client]
command=/usr/bin/python /usr/