# src/main/resources/org/jenkinsci/test/acceptance/docker/fixtures/SvnContainer/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "exec /usr/sbin/apache2 -DFOREGROUND"
